# lottery/redis/crontab/clean.sh
#!/bin/sh


> /home/lau/project/Irean/PHP/redis/crontab/res.txt
> /home/lau/project/Irean/PHP/redis/crontab/res_1.txt

