# lottery/redis/crontab/crontab_1.sh
#!/bin/sh


date >> /home/lau/project/Irean/PHP/redis/crontab/res_1.txt

