# redis/crontab/clean.sh
#!/bin/sh
PWD=`pwd`

> ${PWD}/res.txt
> ${PWD}/res_1.txt

