# redis/crontab/crontab_1.sh
#!/bin/sh

PWD=`pwd`
date >> ${PWD}/res_1.txt

