# redis/crontab/crontab.sh
#!/bin/sh

PWD=`pwd`
date >> ${PWD}/res.txt

